# ops/systemd/ihos-backup.timer
[Unit]
Description=Run the Iron House OS recovery backup daily

[Timer]
OnCalendar=*-*-* 02:15:00 UTC
RandomizedDelaySec=900
Persistent=true
Unit=ihos-backup.service

[Install]
WantedBy=timers.target
