# ops/systemd/ihos-deploy-agent.service
[Unit]
Description=Iron House OS gated production deployment agent
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/iron-house-os/source
Environment=IHOS_REPO_ROOT=/opt/iron-house-os/source
Environment=IHOS_COMPOSE_ENV_FILE=/etc/iron-house-os/production.env
ExecStart=/bin/bash /opt/iron-house-os/source/ops/deploy-agent/deploy.sh
TimeoutStartSec=45min
